# test/assets/auto-recovery/microshift-auto-recovery.service
[Unit]
Description=MICROSHIFT AUTO RECOVERY

[Service]
ExecStart=/usr/bin/echo HELLO

[Install]
WantedBy=multi-user.target